# packages/escalade/libretro/retroarch/system.d/tmp-joypads.mount
[Unit]
Description=Joypad configs directory
Before=graphical.target
After=sysinit.target
DefaultDependencies=no

[Mount]
What=none
Where=/tmp/joypads
Type=overlay
Options=lowerdir=/usr/share/retroarch/joypads,upperdir=/storage/.config/retroarch/joypads,workdir=/storage/.tmp/joypads-workdir

[Install]
WantedBy=multi-user.target
