# packages/oem/filesystem/usr/config/profile.d/README
These scripts will be read on every login. They need to have .conf extension.
