# packages/escalade/network/nfs-utils/system.d/nfs-statedirs.service
[Unit]
Description=Setup NFS state dirs
Before=sysinit.target
After=local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/true
ExecStartPost=-mkdir -p /var/lib/nfs/rpc_pipefs
ExecStartPost=-mkdir -p /var/lib/nfs/v4recovery
ExecStartPost=-mkdir -p /var/lib/nfs/sm
ExecStartPost=-mkdir -p /var/lib/nfs/sm.bak
ExecStartPost=-touch /var/lib/nfs/etab
ExecStartPost=-touch /var/lib/nfs/rmtab
ExecStartPost=-touch /var/lib/nfs/state

[Install]
WantedBy=sysinit.target
