# packages/oem/filesystem/usr/lib/systemd/system/alsa-config.service
[Unit]
Description=Configure ALSA
After=pulseaudio.service
Before=frontend.service

[Service]
Type=oneshot
ExecStart=/usr/bin/alsa-config.sh

[Install]
WantedBy=graphical.target
