# packages/oem/filesystem/usr/lib/systemd/system/frontend.service
[Unit]
Description=Start the default frontend
After=multi-user.target xorg.service

[Service]
Type=oneshot
ExecStart=/usr/bin/frontend.sh
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
